# systemd/magical.service
[Unit]
Description=Calendar update service
Wants=magical.timer

[Service]
User=pi
Group=pi
WorkingDirectory=/home/pi/magical/
Type=oneshot
ExecStart=/home/pi/magical/env/bin/python3 /home/pi/magical/maginkcal.py

[Install]
WantedBy=multi-user.target